# docker/root/etc/services.d/galdir/run
#!/usr/bin/with-contenv bash

echo "Starting GalDir Web Server"

exec \
    waitress-serve --listen=0.0.0.0:6543 --call 'galdir:create_app'
