# docker/root/etc/services.d/galdir/run
#!/usr/bin/with-contenv bash

echo "Starting GalDir Web Server"

exec \
    waitress-serve --listen=${GALDIR_LISTEN} --call 'galdir:create_app'
